# files/99-disable-apple-ir.rules
ACTION=="add", ATTR{idVendor}=="05ac", ATTR{idProduct}=="8242", RUN+="/bin/sh -c 'echo 0 >/sys$DEVPATH/authorized'"
